# MicrosoftAzure/Ansible/roles/docker/templates/docker-compose.conf
[program:docker-compose]
command=docker-compose up
directory={{ docker_compose_home }}
redirect_stderr=true
autostart=true
autorestart=true
priority=10
